# config/imagebuilder/files/etc/uci-defaults/99-passwd
encrypt_password=$(openssl passwd -1 'mywrtadmin')
sed -i -e "s:^root\:[^\:]*\::root\:$encrypt_password\::" /etc/shadow